# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-django/run
#!/command/with-contenv /bin/ash

cd /app/

exec s6-setuidgid memoria granian --interface asginl --loop uvloop --host "127.0.0.1" --port 8100 "memoria.asgi:application"
