# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-valkey/run
#!/command/with-contenv ash
# shellcheck shell=ash

exec s6-setuidgid cleaner \
    /usr/bin/valkey-server \
        /app/valkey/valkey.conf
