# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-huey/run
#!/command/with-contenv /bin/ash

cd /app/

exec s6-setuidgid memoria python3 manage.py run_huey
