# docker/rootfs/etc/s6-overlay/s6-rc.d/init-start/run
#!/command/with-contenv /bin/ash
# shellcheck shell=bash

log_prefix="[init-start]"

echo "${log_prefix} memoria docker container starting..."
