# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/usr/bin/with-contenv bash

# Execute nginx in the foreground
exec s6-setuidgid memoria nginx -g 'daemon off;'
