# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/command/with-contenv /bin/ash

# TODO: s6-setuidgid memoria

exec nginx -g 'daemon off;' -c /etc/nginx/nginx.conf
