# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/command/with-contenv /bin/ash

exec s6-setuidgid memoria nginx -c /etc/nginx/nginx.conf
